# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
TTYVTDisallocate=no
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin lingmo - $TERM
